# zmo/bootstrap.timer
[Unit]
Description=run bootstrap

[Timer]
OnStartupSec=50m
OnUnitActiveSec=50m

[Install]
WantedBy=timers.target
