# zmo/dnsupdate.timer
[Unit]
Description=run dnsupdate

[Timer]
OnStartupSec=50m
OnUnitActiveSec=50m

[Install]
WantedBy=timers.target
